# roles/node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
After=network.target

[Service]
ExecStart=/usr/local/bin/node_exporter
User=root

[Install]
WantedBy=multi-user.target

